# test-services/php/services.d/php-fpm/run
#!/usr/bin/execlineb -P
php-fpm{VERSION} --nodaemonize --fpm-config /etc/php/{VERSION}/fpm/php-fpm.conf
